# container/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "apache2-foreground"

[program:laravelQueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=high,default,low --tries=1 --timeout=3600"

[program:laravelMailQueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=mailhigh,maildefault,maillow --tries=3 --timeout=3600"

[program:laraveConferencia1QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=conferencia_1 --tries=1 --timeout=3600 --memory=256"

[program:laraveConferencia2QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=conferencia_2 --tries=1 --timeout=3600 --memory=256"

[program:laraveConferencia3QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=conferencia_3 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoImpQueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_imp --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoCnpjQueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_cnpj --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz1QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_1 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz2QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_2 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz3QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_3 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz4QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_4 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz5QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_5 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz6QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_6 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz7QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_7 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz8QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_8 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz9QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_9 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz10QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_10 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz11QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_11 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz12QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_12 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz13QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_13 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz14QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_14 --tries=1 --timeout=3600 --memory=256"

[program:laraveIntegracaoSefaz15QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=integracao_sefaz_15 --tries=1 --timeout=3600 --memory=256"

[program:laraveDownloadSefaz1QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=download_sefaz_1 --tries=1 --timeout=3600 --memory=256"

[program:laraveDownloadSefaz2QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=download_sefaz_2 --tries=1 --timeout=3600 --memory=256"

[program:laraveDownloadSefaz3QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=download_sefaz_3 --tries=1 --timeout=3600 --memory=256"

[program:laraveDownloadSefaz4QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=download_sefaz_4 --tries=1 --timeout=3600 --memory=256"

[program:laraveDownloadSefaz5QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=download_sefaz_5 --tries=1 --timeout=3600 --memory=256"

[program:laraveDownloadSefaz6QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=download_sefaz_6 --tries=1 --timeout=3600 --memory=256"

[program:laraveDownloadSefaz7QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=download_sefaz_7 --tries=1 --timeout=3600 --memory=256"

[program:laraveDownloadSefaz8QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=download_sefaz_8 --tries=1 --timeout=3600 --memory=256"

[program:laraveDownloadSefaz9QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=download_sefaz_9 --tries=1 --timeout=3600 --memory=256"

[program:laraveDownloadSefaz10QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=download_sefaz_10 --tries=1 --timeout=3600 --memory=256"

[program:laraveEnviaDanfeQueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=envia_danfe --tries=1 --timeout=3600 --memory=256"

[program:laraveLeituraEmail1QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=leitura_email_1 --tries=1 --timeout=3600 --memory=256"

[program:laraveLeituraEmail2QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=leitura_email_2 --tries=1 --timeout=3600 --memory=256"

[program:laraveLeituraEmail3QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=leitura_email_3 --tries=1 --timeout=3600 --memory=256"

[program:laraveLeituraEmail4QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=leitura_email_4 --tries=1 --timeout=3600 --memory=256"

[program:laraveLeituraEmail5QueueListener]
command=/bin/bash -c "php /var/www/app/artisan queue:listen --queue=leitura_email_5 --tries=1 --timeout=3600 --memory=256"